# services/dnsmasq/run
#!/bin/sh
exec dnsmasq -k -5 -C /etc/dnsmasq.conf
