# services/dnsmasq/run
#!/usr/bin/env sh
exec /usr/sbin/dnsmasq -k
